# custom-systemd-units/evansync.service
[Unit]
Description=Runs a custom cron script

[Service]
ExecStart=/home/myh2910/dotfiles/sh-scripts/cron.sh
